# board/common/overlay/etc/udev/rules.d/90-usb-net.rules
SUBSYSTEM=="net", ACTION=="add", ENV{ID_BUS}=="usb", ENV{INTERFACE}=="eth0", RUN+="/sbin/ifup eth0"
